# recipes-core/kernel-panic/files/sysrq-panic.service
[Unit]
Description=sysrq panic

[Service]
Type=oneshot
ExecStart=/bin/sh -c "echo c > /proc/sysrq-trigger"

[Install]
WantedBy=default.target
