# recipes-core/tmp-fs/files/tmp.mount.tmpl
[Unit]
Description=Create /tmp

[Mount]
What=tmpfs
Where=/tmp
Type=tmpfs
Options=${TMP_FS_OPTIONS}

[Install]
WantedBy=local-fs.target
